# arch/x86_64/tools/irqbalance_mck.service
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
EnvironmentFile=/tmp/irqbalance_mck
ExecStart=/usr/sbin/irqbalance --foreground $IRQBALANCE_ARGS

[Install]
WantedBy=multi-user.target
